# scripts/chat_bridge/caor-chat-bridge-sync.service
[Unit]
Description=Retry ACOR Chat Bridge V3 remote synchronization
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/hanlinxuan/research/ACOR-online-reconstruction
Environment=PATH=/home/hanlinxuan/miniconda3/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash /home/hanlinxuan/research/ACOR-online-reconstruction/scripts/chat_bridge/retry_remote_sync.sh
Nice=10

[Install]
WantedBy=default.target
