# scripts/chat_bridge/caor-chat-bridge-sync.timer
[Unit]
Description=Check ACOR Chat Bridge V3 remote sync every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
